# overseer/over-seer.service
[Unit]
Description=Bergamot Over-Seer — event dashboard server
Documentation=https://github.com/bergamot
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/bergamot/overseer
Environment="TCP_PORT=9000"
Environment="FLASK_HOST=0.0.0.0"
Environment="FLASK_PORT=5000"
ExecStart=/usr/bin/python3 /opt/bergamot/overseer/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=over-seer

[Install]
WantedBy=multi-user.target
