# underseer/under-seer.service
[Unit]
Description=Bergamot Under-Seer — All-Seer polling agent
Documentation=https://github.com/bergamot
After=network.target
# Ensure the All-Seer kernel module is loaded before starting.
# If you use a custom service to load the module, list it here instead:
ConditionPathExists=/proc/all_seer

[Service]
Type=simple
User=root
# Set OVERSEER_HOST to the IP/hostname of the Over-Seer machine.
Environment="OVERSEER_HOST=CHANGE_ME"
Environment="OVERSEER_PORT=9000"
Environment="POLL_INTERVAL_MS=100"
Environment="BATCH_MAX=64"
Environment="RECONNECT_MAX_S=30"
Environment="PROC_PATH=/proc/all_seer"
ExecStart=/usr/bin/python3 /opt/bergamot/underseer/underseer.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=under-seer

[Install]
WantedBy=multi-user.target
